# pi-gen-custom/stage2/10-gps-logger/files/etc/systemd/system/gps-logger.service
[Unit]
Description=GPS Logger Service
After=gpsd.service
Wants=gpsd.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/gps-logger/gps_logger.py
Restart=always
RestartSec=2
User=root
Group=root

[Install]
WantedBy=multi-user.target
